# cc-relay-server.service
[Unit]
Description=Claude Code Relay Server
Documentation=https://github.com/wakaka6/claude-code-relay
After=network.target

[Service]
Type=simple
User=cc-relay
Group=cc-relay
WorkingDirectory=/opt/cc-relay
ExecStart=/opt/cc-relay/cc-relay-server --config /opt/cc-relay/config.toml
Restart=on-failure
RestartSec=5

Environment=RUST_LOG=info

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/cc-relay/data
PrivateTmp=true

[Install]
WantedBy=multi-user.target
